# deploy/systemd.service
# Systemd service file for Perm_shop
# Place this file at: /etc/systemd/system/shoeapp.service
# Then run: sudo systemctl daemon-reload && sudo systemctl enable shoeapp

[Unit]
Description=Shoe Store FastAPI Application
After=network.target

[Service]
User=shoeapp
Group=shoeapp
WorkingDirectory=/home/shoeapp/Perm_shop
Environment="PATH=/home/shoeapp/Perm_shop/.venv/bin"
EnvironmentFile=/home/shoeapp/Perm_shop/.env
ExecStart=/home/shoeapp/Perm_shop/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8002 --workers 2 --log-config /home/shoeapp/Perm_shop/logging.conf
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

